# infra/systemd/dialeros-admin.service
[Unit]
Description=DialerOS admin GUI + control plane (Next.js on port 1111)
Documentation=https://github.com/steven-patrick18/dialeros
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dialeros
Group=dialeros
WorkingDirectory=/opt/dialeros
ExecStart=/usr/bin/pnpm --dir /opt/dialeros --filter @dialeros/admin-gui start
Restart=on-failure
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=PORT=1111
Environment=DIALEROS_DB=/opt/dialeros/apps/admin-gui/data/dialeros.db
Environment=DIALEROS_MASTER_KEY_PATH=/opt/dialeros/apps/admin-gui/data/.master_key
# Iter 145 — load the shared secret + operator-tweakable values
# (KAMAILIO_INBOUND_TOKEN, DIALEROS_REPORT_TO, ...). The leading
# "-" means systemd won't fail if the file is missing (fresh
# installs run bootstrap-admin-env.sh to create it).
EnvironmentFile=-/etc/dialeros/admin.env

# Hardening — minimal. ANY of the systemd "Protect*" / "Restrict*" /
# "Lock*" sandboxing directives implicitly set NoNewPrivileges=yes,
# which blocks sudo from elevating to run the telephony install
# script. Since the install must run from the admin-gui process tree,
# we keep only the lightweight isolation that doesn't trigger
# NoNewPrivileges:
#   - PrivateTmp gives an isolated /tmp without touching privs.
# Everything else relies on:
#   - User=dialeros (the process can't do anything privileged on its own)
#   - /etc/sudoers.d/dialeros-telephony (narrow grant for the one script)
#   - The install script itself being root-owned + 755.
PrivateTmp=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dialeros-admin

[Install]
WantedBy=multi-user.target
