# infra/systemd/dialeros-admin.service
[Unit]
Description=DialerOS admin GUI + control plane (Next.js on port 1111)
Documentation=https://github.com/steven-patrick18/dialeros
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dialeros
Group=dialeros
WorkingDirectory=/opt/dialeros
ExecStart=/usr/bin/pnpm --dir /opt/dialeros --filter @dialeros/admin-gui start
Restart=on-failure
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=PORT=1111
Environment=DIALEROS_DB=/opt/dialeros/apps/admin-gui/data/dialeros.db
Environment=DIALEROS_MASTER_KEY_PATH=/opt/dialeros/apps/admin-gui/data/.master_key

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ReadWritePaths=/opt/dialeros/apps/admin-gui/data /opt/dialeros/apps/admin-gui/.next
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
# AF_NETLINK needed for uv_interface_addresses (Next.js startup network introspection)
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
LockPersonality=yes
RestrictRealtime=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dialeros-admin

[Install]
WantedBy=multi-user.target
