# infra/systemd/dialeros-ai-worker.service
# Iter 137 — local AI worker. Replaces iter-136's OpenAI variant
# with whisper.cpp + Ollama, both running on this host. No
# external network access at runtime.

[Unit]
Description=DialerOS local AI post-call worker (whisper.cpp + Ollama)
After=network.target dialeros-admin.service ollama.service
Wants=ollama.service

[Service]
Type=oneshot
EnvironmentFile=-/etc/dialeros/admin.env
ExecStart=/usr/bin/python3 /opt/dialeros/scripts/ai-worker.py
StandardOutput=journal
StandardError=journal
User=dialeros
Group=dialeros
NoNewPrivileges=true
PrivateTmp=true
# Read access only to recording WAVs + the local model dir;
# write only to a per-run /tmp via PrivateTmp.
ProtectSystem=strict
ReadOnlyPaths=/var/lib/dialeros/recordings /var/lib/dialeros/ai
# whisper.cpp + ollama can chew through CPU on long calls.
# Soft-cap so the worker doesn't starve the rest of the box.
CPUQuota=300%
Nice=10
